# freeswitch/build/supervisor/conf.d/freeswitch.ini
[program:freeswitch]
command=bash /startup_fs.sh
autostart=true
stopsignal=TERM
startretries=100
autorestart=true
stopwaitsecs=5
stopasgroup=true