# freeswitch/build/supervisor/conf.d/memcache.centos.ini
[program:memcached]
command=/usr/bin/memcached -u memcached  -m 256 -p 11211 -l 127.0.0.1  -I 32m
autostart=true
stopsignal=TERM
startretries=100
autorestart=true
stopwaitsecs=5
stopasgroup=true
