# freeswitch/build/supervisor/conf.d/postgresql.ini
[program:postgresql]
#command=/usr/lib/postgresql/13/bin/postgres -D /var/lib/postgresql/13/main -c config_file=/etc/postgresql/13/main/postgresql.conf
user=postgres
autostart=true
stopsignal=TERM
startretries=100
autorestart=true
stopwaitsecs=5
stopasgroup=true